# base_filesystem/var/lib/opkg/info/dropbear.list
/var/
/usr/bin/dbclient
/etc/
/usr/sbin/dropbearmulti
/usr/bin/
/etc/dropbear/
/etc/init.d/dropbear
/usr/sbin/dropbearkey
/etc/init.d/
/usr/sbin/dropbearconvert
/usr/sbin/dropbear
/usr/sbin/
/etc/default/
