# base_filesystem/var/lib/opkg/info/neard.list
/usr/lib/neard/list-adapters
/usr/lib/neard/
/usr/lib/neard/dump-device
/usr/lib/neard/handover-agent
/usr/lib/neard/nfctool
/usr/lib/neard/dump-tag
/usr/lib/neard/enable-adapter
/usr/lib/neard/neard-ui.py
/usr/lib/neard/start-poll
/usr/lib/neard/dump-record
/usr/lib/neard/monitor-near
/usr/lib/neard/write-tag
/usr/lib/neard/disable-adapter
/usr/lib/neard/ndef-agent
/usr/lib/neard/push-device
/usr/lib/neard/stop-poll
/usr/lib/neard/neard
/usr/lib/neard/bt-handover
/etc/init.d/neard
/etc/dbus-1/system.d/org.neard.conf
