# 75XX_TNSTC_FILESYSTEM/var/lib/opkg/info/psplash.list
/mnt/.psplash/
/mnt/
/usr/bin/psplash-write
/etc/init.d/
/etc/init.d/psplash.sh
