# 75XX_TNSTC_FILESYSTEM/var/lib/opkg/info/rpcbind.list
/usr/sbin/rpcbind
/usr/sbin/rpcinfo
/etc/init.d/rpcbind
